# etc/cron.daily/dehydrated
#!/bin/sh -eu

export PATH=/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin
exec >>/var/log/dehydrated.log 2>&1
date
dehydrated-zonefiles-wrapper -c
